# services/udiskie.service
[Unit]
Description=Start udiskie

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/udiskie

[Install]
WantedBy=default.target
